# deploy/launchagents/com.swiftman.mlx-final-text.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.swiftman.mlx-final-text</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/runner/mlx-final-py/scripts/run_mlx_final_text.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/runner/mlx-final-py</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/mlx-final-text.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mlx-final-text.err.log</string>
</dict>
</plist>
